# ansible/roles/pihole/files/systemd/pihole-gravity.service
[Unit]
Description=Pihole Gravity Update service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole-gravity.sh
