# ansible/roles/t2-epg/files/systemd/t2-epg.service
[Unit]
Description=T-2 EPG Download service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/t2-epg.sh
