# ansible/roles/t2-epg/files/systemd/t2-playlist.service
[Unit]
Description=T-2 Playlist Download service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/t2-playlist.sh
